# supervisord.conf
[unix_http_server]
file=/tmp/xblock-supervisor.sock

[supervisord]
logfile=/home/aendra/xblock-docker/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/tmp/xblock-supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/xblock-supervisor.sock

[program:xblock-worker]
command=/home/aendra/.pyenv/versions/xblock/bin/python /home/aendra/xblock-docker/processor/worker.py
directory=/home/aendra/xblock-docker/processor
autostart=true
autorestart=true
startretries=10
startsecs=15
stopwaitsecs=30
stdout_logfile=/home/aendra/xblock-docker/logs/worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
redirect_stderr=true
environment=HOME="/home/aendra",PATH="/home/aendra/.pyenv/versions/xblock/bin:/home/aendra/.pyenv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:xblock-dashboard]
command=/home/aendra/.pyenv/versions/xblock/bin/uvicorn dashboard:app --host 0.0.0.0 --port 8080 --log-level warning --loop asyncio
directory=/home/aendra/xblock-docker/processor
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/home/aendra/xblock-docker/logs/dashboard.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
redirect_stderr=true
environment=HOME="/home/aendra",PATH="/home/aendra/.pyenv/versions/xblock/bin:/home/aendra/.pyenv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
